# Sniffer/install_headless/rc.local
#!/bin/sh -e
sleep 10
x-terminal-emulator -e bash -c 'sudo python C:\studystuff\IoT/sniffer2.0.py headless' &
exit 0
